# ostp.service
[Unit]
Description=ostp

[Service]
ExecStart=/usr/bin/python3 /opt/ostp/ostp.py

[Install]
WantedBy=multi-user.target

